# infra/deploy/systemd/books-api.service
[Unit]
Description=Books Management System -- API
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service

[Service]
Type=simple
User=books
Group=books
WorkingDirectory=/opt/books-manag-system/backend
EnvironmentFile=/etc/books/books.env
Environment=APP_ROLE=api
ExecStart=/opt/books-manag-system/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
