# infra/deploy/systemd/books-beat.service
[Unit]
Description=Books Management System -- Celery Beat
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service

[Service]
Type=simple
User=books
Group=books
WorkingDirectory=/opt/books-manag-system/backend
EnvironmentFile=/etc/books/books.env
Environment=APP_ROLE=beat
ExecStart=/opt/books-manag-system/venv/bin/celery -A app.celery_app.celery_app beat --loglevel INFO
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
