# infra/runtime/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

# Each program logs to stdout/stderr so `docker logs` shows everything,
# prefixed with the program name. autorestart keeps single-process
# crashes from taking the whole stack down.

[program:nginx]
command=nginx -g 'daemon off;'
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=10

[program:api]
command=/usr/local/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
directory=/app
user=books
environment=APP_ROLE="api"
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=20
startsecs=3

[program:worker]
command=/usr/local/bin/celery -A app.celery_app.celery_app worker -Q scan,enrich,maintenance --loglevel INFO
directory=/app
user=books
environment=APP_ROLE="worker"
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=30
startsecs=3

[program:beat]
command=/usr/local/bin/celery -A app.celery_app.celery_app beat --loglevel INFO
directory=/app
user=books
environment=APP_ROLE="beat"
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=30
startsecs=3

[program:admin-web]
command=node /app/admin-web/server.js
directory=/app/admin-web
user=books
environment=PORT="3000",HOSTNAME="127.0.0.1",NODE_ENV="production",NEXT_TELEMETRY_DISABLED="1"
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=40
startsecs=3

[program:reader-web]
command=node /app/reader-web/server.js
directory=/app/reader-web
user=books
environment=PORT="3001",HOSTNAME="127.0.0.1",NODE_ENV="production",NEXT_TELEMETRY_DISABLED="1"
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=40
startsecs=3

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
